# deploy/river-orchestrator.service
[Unit]
Description=River Engine Orchestrator
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=river
Group=river
ExecStart=/usr/local/bin/river-orchestrator --config /etc/river/river.json
EnvironmentFile=/etc/river/river.env
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

[Install]
WantedBy=multi-user.target
